# configs/automatic-backup/nextcloud-deployment-backup.service
[Unit]
Description=NextCloud Deployment Backup

[Service]
ExecStart=/usr/bin/nextcloud-deployment-backup.sh
